# bookshelf_server/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:dart]
command=dart "/app/bin/main.dart"